# deploy/telegram-bot-backup.timer
[Unit]
Description=Run Sepid Telegram Bot verified SQLite backup hourly

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=1min
RandomizedDelaySec=5min
Unit=telegram-bot-backup.service

[Install]
WantedBy=timers.target
